# deploy/supervisord/validator.conf
{% from 'macros' import program with context %}

[supervisord]
minfds=600000

{% if run_validator %}
{{ program('validator', '') }}
command=/home/sol/sol/sol/api.sh
{% endif %}

{% if is_watchtower is defined %}
{{ program('watchtower', '') }}
command=/home/sol/sol/sol/watchtower.sh
{% endif %}

[program:health_check_server]
command=/home/sol/sol/sol/env/bin/python -m health.main
environment=PS=health_check_server,TZ=UTC
directory=/home/sol/sol/sol
startsecs=3
stopwaitsecs=30
user=sol
stopasgroup=true
startretries=100000
autorestart=true
redirect_stderr=true
stdout_logfile_maxbytes=2000000000
stdout_logfile_backups=3
